# scripts/com.fantasywc.commentary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Daily commentary sync. Install with:

    cp scripts/com.fantasywc.commentary.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.fantasywc.commentary.plist

  Runs scripts/sync_commentary.sh every day at 09:00 LOCAL time. If the Mac is
  asleep/off at 09:00, launchd runs it at the next wake. Logs go to /tmp.
  To change the time, edit Hour/Minute below, then:
    launchctl unload ~/Library/LaunchAgents/com.fantasywc.commentary.plist
    cp scripts/com.fantasywc.commentary.plist ~/Library/LaunchAgents/
    launchctl load  ~/Library/LaunchAgents/com.fantasywc.commentary.plist
  To stop it entirely: launchctl unload ~/Library/LaunchAgents/com.fantasywc.commentary.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fantasywc.commentary</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/dwb/Code/fantasy_world_cup/scripts/sync_commentary.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/fwc-commentary.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/fwc-commentary.log</string>
</dict>
</plist>
